# services/wakelock.service
[Unit]
Description=Lock the screen on resume from suspend

[Service]
User=caymard
Type=forking
Environment=DISPLAY=:0
ExecStart=/home/caymard/.config/i3/scripts/lock.sh /home/caymard/.config/i3/wallpapers/space.png

[Install]
WantedBy=sleep.target suspend.target
